# infra/remote/dns-watchdog/dns-watchdog.timer
[Unit]
Description=snoringcat-platform DNS watchdog (hourly)
# Cleans up s-<ip>.<SERVER_DNS_BASE> A records orphaned by
# game-server containers that exited via SIGKILL (no trap
# fired). Records younger than $MAX_RECORD_AGE_HOURS are
# left alone; older ones are deleted.

[Timer]
# Offset 17 minutes past the hour to avoid the round-minute
# pile-up with cost-monitor (00) and other system timers.
OnCalendar=*-*-* *:17:00 UTC
Persistent=true
Unit=dns-watchdog.service
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
