# infra/remote/pg-backup/pg-backup.timer
[Unit]
Description=snoringcat-platform Postgres backup (daily 03:11 UTC)

[Timer]
# 03:11 UTC keeps the backup window clear of cost-monitor's
# 09:00 UTC daily summary, dns-watchdog, and Hetzner's typical
# nightly maintenance windows. Off-round-minute (per the
# scheduling.md offset convention).
OnCalendar=*-*-* 03:11:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
